# install/mailscanner/bin/cron/processing_messages_alert.cron
#!/bin/bash

[ -x /opt/MailScanner/bin/processing_messages_alert ] || exit 0
exec /opt/MailScanner/bin/processing_messages_alert > /dev/null 2>&1
exit 0
