# scripts/systemd/ufw.service.d/override.conf
[Unit]
Description=Uncomplicated firewall for MailCleaner
After=network.service

[Service]
ExecStartPre=/usr/mailcleaner/bin/dump_firewall.pl

[Install]
WantedBy=mailcleaner.target
