# scripts/systemd/ufw.service.d/override.conf
[Unit]
Description=Uncomplicated firewall for MailCleaner

[Service]
ExecStartPre=/usr/mailcleaner/bin/dump_firewall.pl
