# scripts/systemd/mariadb@master-nopass.service.d/override.conf
[Unit]
Conflicts=mariadb@master.service

[Service]
Environment=MYSQLD_OPTS="--skip-grant-tables --port=3306"
