# scripts/systemd/spamd@spamd.service.d/override.conf
[Unit]
Description=SpamAssassin Daemon daemon for MailCleaner
After=mariadb@slave.service
PartOf=mailcleaner.service
Requires=dccifd

[Service]
PIDFile=/var/mailcleaner/run/spamd/%I.pid
User=debian-spamd
Group=mailscanner
ExecStartPre=
ExecStartPre=+/usr/mailcleaner/bin/dump_spamd_config.pl
ExecStart=
ExecStart=/usr/sbin/spamd --socketpath=/var/mailcleaner/spool/spamd/%I.sock --pidfile=/var/mailcleaner/run/spamd/%I.pid --syslog=/var/mailcleaner/log/mailscanner/%I.log --username=debian-spamd --groupname=mailscanner --siteconfigpath=/usr/mailcleaner/share/%I --nouser-config --daemonize
ProtectSystem=false
PrivateTmp=false
PrivateDevices=false
ProtectHostname=false
ProtectClock=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectKernelLogs=false
ProtectControlGroups=false
RestrictNamespaces=false
LockPersonality=false
RestrictRealtime=false
RestrictSUIDSGID=false

[Install]
WantedBy=mailcleaner.target
