# scripts/systemd/mariadb@slave.service.d/port.conf
[Service]
Environment=MYSQLD_OPTS="--port=3307"
