# scripts/systemd/mariadb@slave-nopass.service.d/override.conf
[Unit]
Conflicts=mariadb@slave.service

[Service]
Environment=MYSQLD_OPTS="--skip-grant-tables --port=3307"
